# files/etc/rc.local
#!/bin/sh
# OpenWrt startup script for RPi4 custom build

# Ensure network is up
sleep 5

# Start essential services if not running
/etc/init.d/network restart
/etc/init.d/dnsmasq restart
/etc/init.d/uhttpd restart

# Enable LED indicators (if available)
echo heartbeat > /sys/class/leds/led0/trigger 2>/dev/null || true

# Log successful boot
logger "OpenWrt RPi4 custom build startup completed"

exit 0
